# deploy/supervisor/maediia-worker.conf
; Supervisor config for the ARQ background worker
;
; Deploy to: /etc/supervisor/conf.d/maediia-worker.conf
; Apply:     supervisorctl reread && supervisorctl update
; Start:     supervisorctl start maediia-worker
; Logs:      supervisorctl tail -f maediia-worker

[program:maediia-worker]
command=/projects/maediia_platform/maediia-api/.venv/bin/python -m arq app.worker.WorkerSettings

directory=/projects/maediia_platform/maediia-api
user=www-data
autostart=true
autorestart=true
startretries=5
stopwaitsecs=30
stopsignal=TERM

environment=PYTHONPATH="/projects/maediia_platform/maediia-api"

stdout_logfile=/var/log/supervisor/maediia-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/maediia-worker.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
