# deploy/supervisor/maediia-api.conf
; Supervisor config for the FastAPI application (uvicorn)
;
; Deploy to: /etc/supervisor/conf.d/maediia-api.conf
; Apply:     supervisorctl reread && supervisorctl update
; Start:     supervisorctl start maediia-api
; Restart:   supervisorctl restart maediia-api
; Logs:      supervisorctl tail -f maediia-api

[program:maediia-api]
command=/projects/maediia_platform/maediia-api/.venv/bin/uvicorn app.main:app
        --host 127.0.0.1
        --port 8000
        --workers 4
        --loop uvloop
        --log-level info
        --access-log
        --proxy-headers
        --forwarded-allow-ips=127.0.0.1

directory=/projects/maediia_platform/maediia-api
user=www-data
autostart=true
autorestart=true
startretries=5
stopwaitsecs=30
stopsignal=TERM

; Environment — supplement with values from .env loaded by the app
environment=PYTHONPATH="/projects/maediia_platform/maediia-api"

stdout_logfile=/var/log/supervisor/maediia-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/maediia-api.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
